# etc/supervisor/conf.d/proxy.ini
#
# BEPSwap RPC Proxy
#
[program:rpc-proxy]
command=nginx -g "daemon off;"
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
